# apps/opencode-bot/com.user.telegram.bridge.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.telegram.bridge</string>
    <key>ProgramArguments</key>
    <array>
        <!-- Replace YOUR_USERNAME with your macOS username -->
        <string>/Users/YOUR_USERNAME/.bun/bin/bun</string>
        <string>run</string>
        <string>index.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <!-- Replace with the actual path where you cloned this repo -->
    <string>/Users/YOUR_USERNAME/opencode/telegram-bridge</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/telegram-bridge.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/telegram-bridge.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/YOUR_USERNAME/.bun/bin</string>
    </dict>
</dict>
</plist>
